# dracut/99cryptsetup/coreos-cryptagent.path
[Unit]
Description="Trigger for coreos-cryptagent service"
DefaultDependencies=false
# Order before crypto volumes
Before=initrd-setup-boot.service basic.target cryptsetup.target

[Path]
PathExists=/boot/etc/coreos-cryptagent/
Unit=coreos-cryptagent.service
